# roles/minio_client/templates/etc/profile.d/mc.sh.j2
# {{ ansible_managed }}
# 2023021601

alias mc='mc --config-dir /etc/mc/'
