# roles/glances/templates/etc/profile.d/glances.sh.j2
# {{ ansible_managed }}
# 2022051701

alias top='glances -t 1 --disable-docker'
alias glances='glances -t 1 --disable-docker'
