# roles/glances/templates/etc/profile.d/glances.sh.j2
# {{ ansible_managed }}
# 2026061201

alias top='glances -t 1'
alias glances='glances -t 1'
