# roles/glances/templates/etc/profile.d/glances.sh.j2
# {{ ansible_managed }}
# 2022032401

alias top='glances -t 1'
