# roles/system_update/templates/etc/systemd/system/update-and-reboot.service.j2
# {{ ansible_managed }}
# 2026061001

[Unit]
Description=update-and-reboot Service
# order-only edge (no Wants/Requires): run before the reboot actor so a requested
# reboot waits for this update to finish. Before= an unqueued unit is a no-op. The
# security lane orders itself After= this service (see security-update.service for
# why the weekly lane runs first).
Before=schedule-reboot.service

[Service]
ExecStart=/usr/local/sbin/update-and-reboot
Type=oneshot
User=root
KillMode=process

[Install]
WantedBy=basic.target
