# roles/system_update/templates/etc/systemd/system/update-and-reboot.service.j2
# {{ ansible_managed }}
# 2026060901

[Unit]
Description=update-and-reboot Service
# order-only edges (no Wants/Requires): run after the security transaction so two
# dnf transactions never overlap, and before the reboot actor so the reboot waits
# for this update to finish. After=/Before= a unit that is not queued is a no-op.
After=security-update.service
Before=schedule-reboot.service

[Service]
ExecStart=/usr/local/sbin/update-and-reboot
Type=oneshot
User=root
KillMode=process

[Install]
WantedBy=basic.target
