# roles/podman_containers/templates/etc/systemd/system/podman-auto-update.timer.d/z00-linuxfabrik.conf.j2
# {{ ansible_managed }}
# 2025111201

[Timer]
OnCalendar=
OnCalendar={{ __podman_containers__auto_update_on_calendar }}
