# roles/tools/templates/etc/profile.d/editor.sh
# {{ ansible_managed }}
#2019052401

export EDITOR={{ tools__editor }}
