# roles/tools/templates/etc/profile.d/history.sh.j2
# {{ ansible_managed }}
# 2023102001

HISTSIZE=100000
HISTTIMEFORMAT="%Y-%m-%d %T "
