# roles/objectstore_backup/templates/etc/systemd/system/objectstore-backup.timer.j2
# {{ ansible_managed }}
# 2023021701

[Unit]
Description=Minio ObjectStore Backup Timer

[Timer]
OnCalendar={{ objectstore_backup__on_calendar }}

[Install]
WantedBy=timers.target
