# roles/objectstore_backup/templates/etc/systemd/system/objectstore-backup.service.j2
# {{ ansible_managed }}
# 2023021701

[Unit]
Description=Minio ObjectStore Backup Service

[Service]
ExecStart=/usr/local/bin/mc.sh
Type=oneshot
User=root

[Install]
WantedBy=basic.target
