# roles/stig/templates/etc/systemd/system/aide-check.timer.j2
# {{ ansible_managed }}
# 2021101501

[Unit]
Description=aide-check Timer

[Timer]
OnCalendar={{ stig__aide_calendar }}
Unit=aide-check.service

[Install]
WantedBy=timers.target
