# roles/schedule_reboot/templates/etc/systemd/system/schedule-reboot.timer.j2
# {{ ansible_managed }}
# 2026060401

[Unit]
Description=schedule-reboot Timer

[Timer]
OnCalendar={{ schedule_reboot__on_calendar }}
Unit=schedule-reboot.service

[Install]
WantedBy=timers.target
