# roles/nextcloud/templates/etc/systemd/system/nextcloud-ldap-show-remnants.timer.j2
# {{ ansible_managed }}
# 2022100401

[Unit]
Description=Nextcloud LDAP Show Remnants Service

[Timer]
OnCalendar=monthly
Unit=nextcloud-app-update.service

[Install]
WantedBy=timers.target
