# roles/nextcloud/templates/etc/systemd/system/nextcloud-scan-files.timer.j2
# {{ ansible_managed }}
# 2022101201

[Unit]
Description=Nextcloud Scan Files Timer

[Timer]
OnCalendar={{ nextcloud__on_calendar_scan_files }}
Unit=nextcloud-scan-files.service

[Install]
WantedBy=timers.target
