# roles/nextcloud/templates/etc/systemd/system/nextcloud-jobs.service.j2
# {{ ansible_managed }}
# 2022100401

[Unit]
Description=Nextcloud Background Jobs Service

[Service]
ExecStart=/usr/bin/php --file /var/www/html/nextcloud/cron.php
Type=oneshot
User=apache
