# roles/nextcloud/templates/etc/systemd/system/nextcloud-jobs.service.j2
# {{ ansible_managed }}
# 2022062902

[Unit]
Description=nextcloud background jobs service

[Service]
ExecStart=/usr/bin/php -f /var/www/html/nextcloud/cron.php
Type=oneshot
User=apache
