# roles/icinga_kubernetes/templates/etc/systemd/system/icinga-kubernetes.service.d/kubeconfig.conf.j2
# {{ ansible_managed }}
# 2025020701

[Service]
Environment="KUBECONFIG={{ icinga_kubernetes__kubeconfig_path }}"
