# roles/apache_solr/templates/etc/systemd/system/apache-solr-dump.timer.j2
# {{ ansible_managed }}
# 2026042801

[Unit]
Description=apache-solr-dump Timer

[Timer]
OnCalendar={{ apache_solr__dump_on_calendar }}
Unit=apache-solr-dump.service

[Install]
WantedBy=timers.target
