# roles/mariadb_server/templates/etc/systemd/system/mariadb-dump.timer.j2
# {{ ansible_managed }}
# 2022070601

[Unit]
Description=mariadb-dump Timer

[Timer]
OnCalendar={{ mariadb_server__dump_on_calendar }}
Unit=mariadb-dump.service

[Install]
WantedBy=timers.target
