# roles/wordpress/templates/etc/systemd/system/wordpress-cron.timer.j2
# {{ ansible_managed }}
# 2022103101

[Unit]
Description=Wordpress Cron Timer

[Timer]
OnCalendar=*:0/15

[Install]
WantedBy=timers.target
