# roles/wordpress/templates/etc/systemd/system/wordpress-cron.timer.j2
# {{ ansible_managed }}
# 2022091901

[Unit]
Description=Wordpress Cron Timer for

[Timer]
OnCalendar=*:0/15

[Install]
WantedBy=timers.target
