# roles/chromium_headless/templates/etc/systemd/system/chromium-headless-proxy.service.j2
# {{ ansible_managed }}
# 2026052101

[Unit]
Description=Proxy to on-demand Headless Chromium
Requires=chromium-headless.service
After=chromium-headless.service

[Service]
{# --exit-idle-time was added in systemd 246. EL8 ships systemd 239 and rejects the option, so we omit it there; the proxy (and the Chromium service bound to it) then stays resident once activated. EL9 (systemd 252) and newer support it. #}
{% if ansible_facts['distribution_major_version'] | int >= 9 %}
ExecStart=/usr/lib/systemd/systemd-socket-proxyd --exit-idle-time={{ chromium_headless__idle_timeout }} {{ chromium_headless__listen_address }}:{{ chromium_headless__backend_port }}
{% else %}
ExecStart=/usr/lib/systemd/systemd-socket-proxyd {{ chromium_headless__listen_address }}:{{ chromium_headless__backend_port }}
{% endif %}
PrivateTmp=true
Restart=on-failure
RestartSec=5
