# roles/mongodb/templates/etc/systemd/system/mongodb-dump.service.j2
# {{ ansible_managed }}
# 2022083101

[Unit]
Description=MongoDB Dump Service

[Service]
ExecStart=/usr/local/sbin/mongodb-dump
Type=oneshot
User=root

[Install]
WantedBy=basic.target
