# roles/keycloak/templates/etc/systemd/system/keycloak.service.j2
# {{ ansible_managed }}
# 2023032301

[Unit]
Description=Keycloak Service
After=network.target

[Service]
User=keycloak
Group=keycloak
EnvironmentFile=/etc/sysconfig/keycloak
ExecStart=/opt/keycloak/bin/kc.sh {% if keycloak__mode == 'development' %}start-dev{% else %}start{% endif %}

WorkingDirectory=/opt/keycloak

[Install]
WantedBy=multi-user.target
