# roles/freeipa_server/templates/etc/systemd/system/pki-tomcatd@.service.d/override.conf.j2
# {{ ansible_managed }}
# 2026030301

[Service]
TimeoutStartSec={{ freeipa_server__systemd_timeoutstartsec }}
