# roles/duplicity/templates/etc/systemd/system/duba.service.j2
# {{ ansible_managed }}
# 2022032901

[Unit]
Description=duba (Duplicity Backup) service

[Service]
ExecStart=/usr/local/bin/duba
Type=oneshot
User=root
