# roles/duplicity/templates/etc/systemd/system/duplicity-backup.service.j2
# {{ ansible_managed }}
# 2022021701

[Unit]
Description=duplicity-backup service

[Service]
ExecStart=/usr/local/bin/duplicity-backup
Type=oneshot
User=root
