# roles/icinga2_master/templates/etc/systemd/system/icinga2.service.d/z00-after-sssd.conf.j2
# {{ ansible_managed }}
# 2026070301

# Ensure Icinga 2 starts after SSSD to prevent sudo errors like
# "problem with defaults entries" caused by sudo calls before SSSD is ready.

[Unit]
After=sssd.target
