# roles/icinga2_master/templates/etc/systemd/system/icinga2.service.d/z00-after-sssd.conf.j2
# {{ ansible_managed }}
# 2026070301

[Unit]
After=sssd.target
