# roles/moodle/templates/etc/systemd/system/moodle-cron.timer.j2
# {{ ansible_managed }}
# 2024072501

[Unit]
Description=Moodle cron.php Timer

[Timer]
OnCalendar={{ moodle__on_calendar_cron }}
Unit=moodle-cron.service

[Install]
WantedBy=timers.target
