# roles/apache_tomcat/templates/etc/systemd/system/tomcat.service.d/override.conf.j2
[Service]
Environment="CATALINA_OPTS=-Xms{{ apache_tomcat__env_xms | d('1024M') }} -Xmx{{ apache_tomcat__env_xmx | d('1024M') }} -server -XX:{{ apache_tomcat__env_xx | d('+UseParallelGC') }}"
