# roles/logstash/templates/etc/systemd/system/logstash.service.d/z10-config-reload.conf.j2
# {{ ansible_managed }}
# 2026012701

[Service]
ExecReload=/bin/kill -HUP $MAINPID
